# edge-devices/bootc/root/etc/systemd/system/flightctl-agent.service.d/override.conf
[Service]
ExecStartPre=/usr/local/bin/update-config.sh
